# contrib/init/wisprd.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env WISPRD_BIN="/usr/bin/wisprd"
env WISPRD_USER="wispr"
env WISPRD_GROUP="wispr"
env WISPRD_PIDDIR="/var/run/wisprd"
# upstart can't handle variables constructed with other variables
env WISPRD_PIDFILE="/var/run/wisprd/wisprd.pid"
env WISPRD_CONFIGFILE="/etc/wispr/wispr.conf"
env WISPRD_DATADIR="/var/lib/wisprd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # wisprd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$WISPRD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run wisprd."
        echo "The setting must appear in $WISPRD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$WISPRD_PIDDIR"
    chmod 0755 "$WISPRD_PIDDIR"
    chown $WISPRD_USER:$WISPRD_GROUP "$WISPRD_PIDDIR"
    chown $WISPRD_USER:$WISPRD_GROUP "$WISPRD_CONFIGFILE"
    chmod 0660 "$WISPRD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$WISPRD_PIDFILE" \
    --chuid $WISPRD_USER:$WISPRD_GROUP \
    --exec "$WISPRD_BIN" \
    -- \
    -pid="$WISPRD_PIDFILE" \
    -conf="$WISPRD_CONFIGFILE" \
    -datadir="$WISPRD_DATADIR" \
    -disablewallet \
    -daemon

